# scripts/launchd/com.boatrace.ai.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.boatrace.ai.morning</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/yu01/projects/boatrace-ai/scripts/morning.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>30</integer>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/yu01/.boatrace-ai/logs/launchd_morning.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/yu01/.boatrace-ai/logs/launchd_morning.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/yu01/.pyenv/shims:/Users/yu01/.pyenv/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
</dict>
</plist>
